# deploy/wahlumfragen-backup.service
# One-shot unit that backs up the Postgres database. Triggered by
# wahlumfragen-backup.timer — do NOT enable this unit directly; enable the
# .timer instead.
#
# Adjust the placeholders for your server before installing:
#   User=             a non-root system user that can run `docker exec`
#                     (in the `docker` group), owning the deploy dir
#   WorkingDirectory= absolute path to the checked-out repo
#   ReadWritePaths=   the backup directory (must match BACKUP_DIR)
#
# Backup tuning (BACKUP_DIR, BACKUP_KEEP_DAYS, …) can be set in the repo .env;
# it is loaded here so the same file configures the whole deployment.

[Unit]
Description=wahlumfragen Postgres backup (one-shot)
After=network-online.target docker.service
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
User=wahlumfragen
WorkingDirectory=/opt/wahlumfragen
EnvironmentFile=-/opt/wahlumfragen/.env
ExecStart=/opt/wahlumfragen/deploy/backup-db.sh

TimeoutStartSec=600

# Hardening — needs to write only the backup directory.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ReadWritePaths=/opt/wahlumfragen/backups
